# assets/supervisord.conf
[program:nginx]
startsecs = 0
autorestart = false
command=/usr/sbin/nginx -g "daemon off;"
